# systemd/reboot.timer
[Unit]
Description = Weekly Reboot Timer

[Timer]
OnCalendar = Mon 3:00 Asia/Tokyo

[Install]
WantedBy = timers.target
