# systemd/reboot.service
[Unit]
Description = Reboot
RefuseManualStart = true
RefuseManualStop = true

[Service]
ExecStart = /sbin/reboot
